# aws-files/startup.service
[Unit]
Description=Startup Script

[Service]
ExecStart=/home/ec2-user/snake/aws-files/run-snake-server.sh

[Install]
WantedBy=multi-user.target
