# aws-files/startup.service
[Unit]
Description=Startup Script

[Service]
ExecStart=/usr/local/bin/run-snake-server.sh

[Install]
WantedBy=multi-user.target
